# tools/launchd/com.rnp.uvicorn.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rnp.uvicorn.service</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/oscaralvarado/venv/bin/uvicorn</string>
    <string>main:app</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8000</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/oscaralvarado/Documents/New project</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>MPLCONFIGDIR</key>
    <string>/tmp/matplotlib</string>
    <key>AUTH_ENABLED</key>
    <string>false</string>
    <key>FORCE_HTTPS</key>
    <string>false</string>
    <key>ENABLE_HSTS</key>
    <string>false</string>
    <key>SECURE_COOKIES</key>
    <string>false</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/rnp_uvicorn.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/rnp_uvicorn.err.log</string>
</dict>
</plist>
